# lib/wanda/wanda.service
[Unit]
Description=Wanda AI Assistant for Scribble
After=network.target

[Service]
Type=simple
User=isdood
ExecStart=/usr/local/bin/wanda-service
Environment=WANDA_CONFIG=/etc/wanda/config.toml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
