# lib/wanda/support/wanda.service
[Unit]
Description=Wanda AI Assistant Service
After=network.target

[Service]
Type=simple
Environment=RUST_LOG=info
ExecStart=%h/.local/bin/wandad
WorkingDirectory=%h
RuntimeDirectory=wanda
LogsDirectory=wanda
ConfigurationDirectory=wanda

# Ensure directories exist
ExecStartPre=/bin/mkdir -p %h/.local/share/wanda
ExecStartPre=/bin/mkdir -p %h/.cache/wanda
ExecStartPre=/bin/mkdir -p %h/.config/wanda

[Install]
WantedBy=default.target
